# home/user/.local/crontab/low-battery-notification
* * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /bin/bash ~/.local/bin/low-battery-notification.sh
